# root/lib/systemd/system/clamd@rspamd.service
.include /lib/systemd/system/clamd@.service

[Unit]
Description = Generic clamav rspamd daemon

[Install]
WantedBy = multi-user.target
